# profiles/zeus/uci-defaults/99-firewall-config
(
sleep 120
uci set firewall.defaults.forward='REJECT'
uci set firewall.defaults.input='ACCEPT'
uci set firewall.defaults.output='ACCEPT'
uci set firewall.defaults.synflood_protect='1'
uci set firewall.defaults.flow_offloading='1'
uci set firewall.defaults.flow_offloading_hw='1'
uci commit firewall
/etc/init.d/firewall restart
) &
exit 0
